# components/puppet/modules/profile/templates/trac_fcgi.service.erb
# -- DO NOT EDIT -- Generated by puppet
# <%= file %> 
#
# systemd.service(5) file for fcgiwrap

[Unit]
Description=Fast cgi service for trac
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
Type=simple

User=nginx
Group=nginx
UMask=006

ExecStartPre=/usr/bin/install -d -o nginx -g nginx -m 775 <%= @trac_fcgi_skdir %>

ExecStart=/usr/local/sbin/trac_fcgi.py '<%= @trac_fcgi_socket %>'
