# components/puppet/modules/profile/templates/tracd.service.erb
# -- DO NOT EDIT -- Generated by puppet
# <%= file %> 
#
# systemd.service(5) file for trac_fcgi.py

[Unit]
Description=trac wiki ticketing system
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
Type=simple

User=trac
Group=nginx
UMask=006

ExecStart= /usr/sbin/tracd --protocol=fcgi --env-parent-dir=/home/trac/products --base-path=/trac --auth=*,/home/trac/.tracdigest,smk --hostname=<%= @tracd_socket_fn %>
