# components/puppet/modules/profile/templates/trac_fcgi.service.erb
# -- DO NOT EDIT -- Generated by puppet
# <%= file %> 
#
# systemd.service(5) file for trac_fcgi.py

[Unit]
Description=Fast cgi service for trac
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
Type=simple

User=trac
Group=nginx
UMask=006

ExecStartPre=/usr/bin/install -d -o nginx -g nginx -m 775 <%= @trac_fcgi_skdir %>
ExecStartPre=/usr/bin/install -d -o trac -g nginx -m 775 /var/log/trac

ExecStart= /usr/sbin/tracd --protocol=fcgi --env-parent-dir=/home/trac/products --base-path=/trac --auth=*,/home/trac/.tracdigest,smk --hostname=<%= @trac_fcgi_socket_fn %>
